# files/mailhog/mailhog.service
[Unit]
Description=MailHog service

[Service]
ExecStart= ~/go/bin/MailHog

[Install]
WantedBy=multi-user.target