# files/mailhog/mailhog.service
[Unit]
Description=MailHog service

[Service]
ExecStart=~/go/bin/Mailhog

[Install]
WantedBy=multi-user.target